# Docker-Images/Ubuntu/CPU/Butterfly/.config/butterfly.conf
[program:butterfly]
command=/start.sh
stderr_logfile=/headless/.config/butterfly.log
stdout_logfile=/headless/.config/butterfly.log
autostart=true
