# Docker-Images/Ubuntu/CPU/Butterfly/.config/butterfly.conf
[program:butterfly]
command=butterfly.server.py --host=biobox --port=5757 --login
stderr_logfile=/headless/.config/butterfly.log
stdout_logfile=/headless/.config/butterfly.log
autostart=true
