# Docker-Images/Ubuntu/CPU/Jupyter/.config/jupyter.conf
[program:JupyterNotebook]
command=start-notebook.sh
stderr_logfile=/workdir/jupyter_err.log
stdout_logfile=/workdir/jupyter.log
autostart=true
