# Docker-Images/Ubuntu/CPU/BioShell/.config/bioshell.conf
[program:BioShell]
command=/bioshell.sh
stderr_logfile=/.config/bioshell.log
stdout_logfile=/.config/bioshell.log
autostart=true
