# Docker-Images/Ubuntu/CPU/BioGUI/.config/biogui.conf
[program:biogui]
command=python3 /headless/BioGUI/launcher.py
stderr_logfile=/headless/BioGUI/guiserver.log
stdout_logfile=/headless/BioGUI/guiserver.log
autostart=true
