# Docker-Images/Ubuntu/GPU/BioGUI/CUDNN/.config/biogui.conf
[program:biogui]
command=python3 /BioGUI/launcher.py
stderr_logfile=/BioGUI/guiserver.log
stdout_logfile=/BioGUI/guiserver.log
autostart=true
