# Docker-Images/Ubuntu/GPU/BioShell/CUDNN/.config/butterfly.conf
[program:BioBash]
command=/butterfly.sh
stderr_logfile=/.config/butterfly.log
stdout_logfile=/.config/butterfly.log
autostart=true
